# snapcast-client/rootfs/etc/s6-overlay/s6-rc.d/snapclient/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: Snapcast Client
# Runs the snapclient
# ==============================================================================
#!/usr/bin/env bashio

declare -a options

bashio::log.info "Starting SnapClient..."

if bashio::config.has_value 'host'; then
    options+=(-h "$(bashio::config 'host')")
fi
if bashio::config.has_value 'port'; then
    options+=(-p "$(bashio::config 'port')")
fi
if bashio::config.has_value 'sampleformat'; then
    options+=(--sampleformat "$(bashio::config 'sampleformat')")
fi
exec snapclient "${options[@]}"